# infra/systemd/sunrise-studio.service
[Unit]
Description=Sunrise Studio Node.js application
After=network.target

[Service]
Type=simple
User=varun_lp
Group=varun_lp
WorkingDirectory=/home/varun_lp/Sunrise_Studio
ExecStart=/usr/local/nvm/versions/node/v24.15.0/bin/npm start
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
